# projects/project-2.service
[Unit]
Description=Project 2
After=networking.service
[Service]
ExecStart=/usr/bin/project-2_server 12345
[Install]
WantedBy=multi-user.target
