# projects/project-2.service
[Unit]
Description=
After=networking.service
[Service]
ExecStart=
[Install]
WantedBy=multi-user.target
